# archive2026_03/rockpis/led-matrix.service
[Unit]
Description=LED Matrix Controller (Q-SYS) — RADXA Rock Pi S
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/led-matrix
ExecStart=/usr/bin/python3 /opt/led-matrix/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Uncomment to run in NO_DISPLAY test mode (no bus error, but no physical display)
# Useful for testing UDP/Web functionality without hardware
# Environment="LED_MATRIX_NO_DISPLAY=1"

# Allow the process to access /dev/mem and GPIO
PrivateDevices=no
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
